# systemd/awww-daemon.service.template
[Unit]
Description=AWWW wallpaper daemon
Documentation=https://codeberg.org/LGFae/awww
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/awww-daemon
Restart=on-failure
RestartSec=3
Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=default.target