# systemd/wallflow.service
[Unit]
Description=Wallflow Wallpaper Daemon
Documentation=https://github.com/MindkraftStudiosGroup/wallflow
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/wallflow daemon start --foreground
Restart=on-failure
RestartSec=5

# Pass graphical session environment variables for display detection
PassEnvironment=WAYLAND_DISPLAY XDG_SESSION_TYPE DISPLAY XDG_CURRENT_DESKTOP DBUS_SESSION_BUS_ADDRESS

# Ensure ~/.local/bin and ~/.cargo/bin are in PATH for backend tools (awww, swww, etc.)
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=graphical-session.target
